# GPS/configs/adagps.conf
description "ADA GPS startup script"
author "Paul"
 
start on runlevel [2345]
stop on runlevel [016]
chdir /home/pi/Desktop/GPS/
exec python /home/pi/Desktop/GPS/run.py
respawn